# modules/systemd-user/payload/rss-refresh.service
[Unit]
Description=RSS Reloader Helper

[Service]
Type=oneshot
ExecStart=/usr/bin/newsboat -x reload
